# supervisor.ini
; Supervisor config template for Uberspace.
;
; Copy to ~/etc/services.d/energychief.ini on Uberspace,
; then run: supervisorctl reread && supervisorctl update
;
; CWD is set to ~/app/, so .env is found on import by shared/config.py

[program:energychief]
command=%(ENV_HOME)s/app/.venv/bin/python %(ENV_HOME)s/app/serve.py
directory=%(ENV_HOME)s/app
autostart=true
autorestart=true
stderr_logfile=%(ENV_HOME)s/app/logs/supervisor-stderr.log
stdout_logfile=%(ENV_HOME)s/app/logs/supervisor-stdout.log
environment=HOME="%(ENV_HOME)s"