# supervisor.ini
; Supervisor config for Uberspace.
;
; Copy to ~/etc/services.d/Enchief.ini on Uberspace,
; then run: supervisorctl reread && supervisorctl update

[program:Enchief]
directory=%(ENV_HOME)s/apps/Energychief
command=%(ENV_HOME)s/apps/Energychief/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000 --no-access-log
autostart=yes
autorestart=yes
startsecs=5
stderr_logfile=%(ENV_HOME)s/apps/Energychief/logs/uvicorn.err.log
stdout_logfile=%(ENV_HOME)s/apps/Energychief/logs/uvicorn.out.log